# signage-kiosk.service
[Unit]
Description=Chromium Kiosk für Signage Anzeige
After=network-online.target signage.service
Wants=network-online.target
Requires=signage.service

[Service]
Type=simple
User=root
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/wag/.Xauthority
WorkingDirectory=/home/wag/signage
ExecStart=/home/wag/signage/scripts/kiosk-start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
